# config/supervisor/nextjs.conf
# House of Veritas - Next.js Supervisor Configuration
# Place this file at: /etc/supervisor/conf.d/nextjs.conf

[program:nextjs]
command=yarn start
environment=HOST="0.0.0.0",PORT="3000",NODE_ENV="production"
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nextjs.err.log
stdout_logfile=/var/log/supervisor/nextjs.out.log
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

# Build process runs before start (optional - run manually or via CI)
# To rebuild: cd /app && yarn build

# Health check endpoint: http://localhost:3000/api/stats
